# algorithm/supervisord.conf
[supervisord]
nodaemon=true
env=WORKSPACE_DIR=/home/workspace/AskOnce/algorithm

[program:answer_by_documents]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key %(ENV_LLM_API_KEY)s --model_name %(ENV_LLM_MODEL_NAME)s --platform_api_url %(ENV_LLM_API_URL)s --search_url %(ENV_SEARCH_URL)s"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:data_convert]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s  --correct_dict_path /home/workspace/AskOnce/assets/dict.json"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:generate_outlines]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key %(ENV_LLM_API_KEY)s --model_name %(ENV_LLM_MODEL_NAME)s --platform_api_url %(ENV_LLM_API_URL)s"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:result_add_reference]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --stop_word_file /home/workspace/AskOnce/assets/stopwords.txt"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:split_question]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key %(ENV_LLM_API_KEY)s --model_name %(ENV_LLM_MODEL_NAME)s --platform_api_url %(ENV_LLM_API_URL)s"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:document_split]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s "
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:question_rewrite]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key %(ENV_LLM_API_KEY)s --model_name %(ENV_LLM_MODEL_NAME)s --platform_api_url %(ENV_LLM_API_URL)s"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:net_rag_assessment]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key %(ENV_LLM_API_KEY)s --model_name %(ENV_LLM_MODEL_NAME)s --platform_api_url %(ENV_LLM_API_URL)s"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[program:search_result_post_process]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

