# algorithm/supervisord.conf
[supervisord]
nodaemon=true
env=WORKSPACE_DIR=/home/workspace/AskOnce/algorithm

[program:answer_by_documents]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key ${LLM_API_KEY} --model_name ${LLM_MODEL_NAME} --platform_api_url ${LLM_API_URL} --search_url ${SEARCH_URL}"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log

[program:data_convert]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s  --correct_dict_path /home/workspace/AskOnce/assets/dict.json"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log


[program:generate_outlines]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key ${LLM_API_KEY} --model_name ${LLM_MODEL_NAME} --platform_api_url ${LLM_API_URL}"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log


[program:result_add_reference]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --stop_word_file /home/workspace/AskOnce/assets/stopwords.txt"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log


[program:split_question]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s  --api_key ${LLM_API_KEY} --model_name ${LLM_MODEL_NAME} --platform_api_url ${LLM_API_URL}"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log

[program:document_split]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s "
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log


[program:question_rewrite]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s  --api_key ${LLM_API_KEY} --model_name ${LLM_MODEL_NAME} --platform_api_url ${LLM_API_URL}"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log


[program:net_rag_assessment]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s --api_key ${LLM_API_KEY} --model_name ${LLM_MODEL_NAME} --platform_api_url ${LLM_API_URL}"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log


[program:search_result_post_process]
command=/bin/bash -c "%(ENV_WORKSPACE_DIR)s/%(program_name)s.bin --tasktype %(program_name)s --worker_name %(program_name)s_%(process_num)s"
process_name=%(program_name)s_%(process_num)s
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/%(program_name)s_%(process_num)s.err.log
stdout_logfile=/var/log/%(program_name)s_%(process_num)s.out.log

